# etc/xkcd.service
[Unit]
Description=Download last xkcd picture

[Service]
Type=simple
User=mcanal
ExecStart=/home/mcanal/sh_script/get_image.sh